# packages/chik-blockchain/chik_blockchain-2.4.2.tar.gz/chik_blockchain-2.4.2/build_scripts/assets/systemd/chik-introducer@.service
[Unit]
Description=Chik Introducer Service for %i
Requires=chik-daemon@%i.service
After=chik-daemon@%i.service

[Service]
Type=simple
Environment=CHIK_ROOT=/home/%i/.chik/mainnet
ExecStart=/opt/chik/start_introducer
User=%i
Group=%i
LimitNOFILE=1048576
LimitNPROC=1048576

[Install]
WantedBy=multi-user.target
