# packages/chik-blockchain/chik_blockchain-2.4.0.tar.gz/chik_blockchain-2.4.0/build_scripts/assets/systemd/chik-wallet@.service
[Unit]
Description=Chik Wallet Service for %i
Requires=chik-daemon@%i.service
After=chik-daemon@%i.service

[Service]
Type=simple
Environment=CHIK_ROOT=/home/%i/.chik/mainnet
ExecStart=/opt/chik/start_wallet
ExecStartPost=/bin/bash -c '(while ! nc -z -v -w1 localhost 10490 2>/dev/null; do echo "Waiting for the wallet to listen on port 10490..."; sleep 1; done); sleep 1'
User=%i
Group=%i
LimitNOFILE=1048576
LimitNPROC=1048576

[Install]
WantedBy=multi-user.target
